# moonshot-current-screen.timer
[Unit]
Description=Refresh current moonshot trading workbench every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=moonshot-current-screen.service

[Install]
WantedBy=timers.target
